# applications/luci-app-nginx/root/etc/uci-defaults/luci-nginx
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@nginx[-1]
	add ucitrack nginx
	set ucitrack.@nginx[-1].init=nginx
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
