# applications/luci-app-dnscrypt-proxy2/root/etc/uci-defaults/luci-dnscrypt-proxy2
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@dnscrypt-proxy2[-1]
	add ucitrack dnscrypt-proxy2
	set ucitrack.@dnscrypt-proxy2[-1].init=dnscrypt-proxy
	commit ucitrack
EOF

exit 0
