# applications/luci-app-frps/root/etc/uci-defaults/luci-frps
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@frps[-1]
	add ucitrack frps
	set ucitrack.@frps[-1].init=frps
	commit ucitrack
EOF

exit 0
