# applications/luci-app-dockerman/root/etc/uci-defaults/luci-app-dockerman
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@dockerman[-1]
	add ucitrack dockerman
	set ucitrack.@dockerman[-1].init=dockerman
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache

exit 0
