# applications/luci-app-macvlan/root/etc/hotplug.d/iface/01-macvlan
#!/bin/sh

[ "$ACTION" = "ifup" ] || exit 0
case "$DEVICE" in
	vth*)
		exit 0
		;;
	*)
		/etc/init.d/macvlan start
		;;
esac
