# applications/luci-app-mariadb/root/etc/uci-defaults/luci-mysqld
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@mysqld[-1]
	add ucitrack mysqld
	set ucitrack.@mysqld[-1].init=mysqld
	commit ucitrack
EOF

exit 0
